# frothiq-control-center.service
[Unit]
Description=FrothIQ Control Center Backend
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=frappe
WorkingDirectory=/home/frappe/frothiq-control-center
ExecStart=/home/frappe/env/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8002 \
    --workers 2 \
    --loop uvloop \
    --http h11 \
    --access-log \
    --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=frothiq-control-center

[Install]
WantedBy=multi-user.target
